# deploy/launchd/com.d.mcp-trust-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template. Placeholders (__REPO_ROOT__, __LOG_OUT__, __LOG_ERR__, __PATH__) are
  substituted by deploy/launchd/install.sh at install time, so this file stays
  machine-agnostic and committable. Do not edit the installed copy here; re-run
  install.sh instead.

  This job is refresh/build only and has no production deployment authority.
  Installation leaves the label persistently disabled. A separate operator
  decision is required before any future scheduling enablement.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.d.mcp-trust-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/refresh_and_publish.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__LOG_OUT__</string>
  <key>StandardErrorPath</key>
  <string>__LOG_ERR__</string>
  <key>RunAtLoad</key>
  <false/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>MCP_TRUST_SITE_BASE_URL</key>
    <string>https://mcp-trust.vercel.app</string>
    <key>MCP_TRUST_SCHEDULER_CONTEXT</key>
    <string>1</string>
  </dict>
</dict>
</plist>
